# quinn/systemd/bubble-packet-filter.service
# Reapply packet-level nftables filter on boot (after network is up).
# Install: copy to /etc/systemd/system/, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable bubble-packet-filter
# Ensure .env with PACKET_MAX_BYTES and WG_INTERFACE is in BUBBLE_ROOT (e.g. /opt/bubble).
# Requires root to run nft; run as root or use sudo in ExecStart.

[Unit]
Description=Bubble packet filter (nftables on WireGuard)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=BUBBLE_ROOT=/opt/bubble
EnvironmentFile=-/etc/default/bubble
# Source .env and run enable; requires root.
ExecStart=/bin/sh -c '"${BUBBLE_ROOT}/scripts/setup_packet_filter.sh" enable || true'

[Install]
WantedBy=multi-user.target
